# systemd/framecast-health.timer
[Unit]
Description=FrameCast Post-Update Health Check Timer

[Timer]
OnBootSec=90

[Install]
WantedBy=timers.target
