# systemd/framecast-update.service
[Unit]
Description=FrameCast OTA Update Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/framecast
ExecStart=/opt/framecast/scripts/update-check.sh
TimeoutStartSec=300
ProtectSystem=strict
NoNewPrivileges=true
ReadWritePaths=/opt/framecast /var/lib/framecast
StandardOutput=journal
StandardError=journal
