# systemd/framecast-kiosk.service
[Unit]
Description=FrameCast Kiosk Display
After=framecast.service graphical.target
Requires=framecast.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/framecast/kiosk/kiosk.sh
Restart=always
RestartSec=10
WatchdogSec=60
MemoryMax=512M

[Install]
WantedBy=graphical.target
