# systemd/framecast-health.service
[Unit]
Description=FrameCast Post-Update Health Check
After=framecast.service framecast-kiosk.service
Wants=framecast.service framecast-kiosk.service

[Service]
Type=oneshot
ExecStart=/opt/framecast/scripts/health-check.sh
TimeoutStartSec=180
ProtectSystem=strict
ReadWritePaths=/opt/framecast /var/lib/framecast
NoNewPrivileges=true
ProtectHome=read-only
ProtectKernelTunables=true
SyslogIdentifier=framecast-health
StandardOutput=journal
StandardError=journal
