# systemd/framecast-schedule.timer
[Unit]
Description=FrameCast Display Schedule Check

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target
