# systemd/framecast.service
[Unit]
Description=FrameCast Web App
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/framecast/app
ExecStart=/usr/bin/gunicorn -c /opt/framecast/app/gunicorn.conf.py web_upload:app
Restart=always
RestartSec=5
MemoryMax=512M
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/framecast /var/lib/framecast /home/pi/media /home/pi/framecast
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
