# systemd/framecast-kiosk.service
[Unit]
Description=FrameCast Kiosk Display
After=framecast.service graphical.target local-fs.target
BindsTo=framecast.service
StartLimitIntervalSec=600
StartLimitBurst=10
StartLimitAction=reboot

[Service]
Type=simple
User=pi
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/framecast/kiosk/kiosk.sh
Restart=always
RestartSec=10
MemoryMax=512M
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/pi/framecast

[Install]
WantedBy=graphical.target
