# systemd/framecast-kiosk.service
[Unit]
Description=FrameCast Kiosk Display
After=framecast.service local-fs.target
BindsTo=framecast.service
StartLimitIntervalSec=600
StartLimitBurst=10
StartLimitAction=none

[Service]
Type=simple
User=pi
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/framecast/kiosk/kiosk.sh
Restart=always
RestartSec=10
MemoryMax=512M
NoNewPrivileges=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ReadWritePaths=/opt/framecast
SyslogIdentifier=framecast-kiosk

[Install]
WantedBy=multi-user.target
