# systemd/photo-upload.service
[Unit]
Description=Pi Photo Display - Web Upload Server
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/pi-photo-display/app
ExecStart=/usr/bin/python3 /opt/pi-photo-display/app/web_upload.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
