# systemd/wifi-manager.service
[Unit]
Description=FrameCast WiFi Manager
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/framecast/scripts/wifi-check.sh
ProtectSystem=strict
NoNewPrivileges=true
ProtectHome=yes
PrivateTmp=true
ReadWritePaths=/var/lib/framecast
SyslogIdentifier=wifi-manager
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
